# deploy/scheduler.crontab
# Add to the deploy user's crontab: crontab -e
* * * * * cd /var/www/crm-system && php artisan schedule:run >> /dev/null 2>&1
